# worker/com.calldeck.worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.calldeck.worker</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>-u</string>
        <string>-m</string>
        <string>worker.worker</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/thegeshwar/calldeck</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/thegeshwar/calldeck/worker/logs/worker_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/thegeshwar/calldeck/worker/logs/worker_stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/thegeshwar/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>CALLDECK_URL</key>
        <string>https://calldeck.thegeshwar.com</string>
        <key>CALLDECK_WORKER_KEY</key>
        <string>dbdc39eba733f19bea4def929b89d6513c7501df586050b0</string>
        <key>SUPABASE_URL</key>
        <string>https://api.qms.thegeshwar.com</string>
        <key>SUPABASE_SERVICE_KEY</key>
        <string>eyJhbGciOiJIUzI1NiIsInR5cCI6IkpXVCJ9.eyJpc3MiOiJzdXBhYmFzZS1kZW1vIiwicm9sZSI6InNlcnZpY2Vfcm9sZSIsImV4cCI6MTk4MzgxMjk5Nn0.EGIM96RAZx35lJzdJsyH-qQwv8Hdp7fsn3W0YpN81IU</string>
    </dict>
</dict>
</plist>
